# ip_rotator/scripts/rotate.cron
*/30 * * * * /bin/bash /home/rotate-ip.sh
@reboot /bin/bash /home/rotate-ip.sh

